# meta-t0-crs/recipes-core/crs-base-config/files/configfs.mount
[Unit]
Description=configfs
DefaultDependencies=no

[Mount]
What=configfs
Where=/configfs
Type=configfs

[Install]
WantedBy=multi-user.target
