# meta-t0-crs/recipes-core/crs-base-config/files/home.mount
[Unit]
Description=Mount /home data partition
DefaultDependencies=no
Before=local-fs.target

[Mount]
What=LABEL=home
Where=/home
Type=ext4
Options=defaults,noatime

[Install]
WantedBy=local-fs.target
